# install/etc/cron.d/button
@reboot pi cd /home/pi/buttons/mk4 && ./buttons > /home/pi/buttons.log
