# install/etc/cron.d/button
@reboot pi cd /home/pi/pibuttons && ./buttons > /home/pi/pibuttons/run.log
