# install/etc/cron.d/self-update
@reboot pi cd /home/pi/pibuttons && ./git pull
